# images/haproxy-1.8/ROOT/etc/services.d/discovery/run
#!/usr/bin/with-contenv /bin/sh

set -e 

echo "Starting background discovery..."

/srv/bg-discovery.sh